# conf/supervisor/mysql.conf
[program:mysql]
command=service mysql start
user=root
autostart=true
autorestart=true
startsecs=3
