# conf/supervisor/openresty.conf
[program:openrestry]
command=/usr/local/openresty/nginx/sbin/nginx
user=root
autostart=true
autorestart=true
startsecs=3
