# conf/supervisor/redis.conf
[program:redis]
command=/usr/bin/redis-server /etc/redis/redis.conf
user=root
autostart=true
autorestart=true
startsecs=3
